# root/etc/services.d/papermerge/run
#!/usr/bin/with-contenv bash

cd /app/papermerge || exit

exec \
	s6-setuidgid abc /usr/bin/uwsgi uwsgi.ini
